# roles/glances/templates/glancesweb.service.j2
[Unit]
Description=GlancesWeb
Requires=systemd-networkd.service
After=systemd-networkd.service

[Timer]
OnBootSec=10

[Service]
ExecStart=/usr/bin/glances -w

[Install]
WantedBy=multi-user.target
