# roles/glances/templates/glancesweb.service.j2
[Unit]
Description=GlancesWeb
After=systemd-networkd.service

[Service]
ExecStart=/usr/bin/glances -w

[Install]
WantedBy=multi-user.target
